# systemd/cerebellum-containers.service
[Unit]
Description=Cerebellum containers (ollama, qdrant, neo4j)
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/psyched/containers/cerebellum
Environment=COMPOSE_PROJECT_NAME=cerebellum
ExecStart=/usr/bin/docker compose -f /opt/psyched/containers/cerebellum/docker-compose.yml up -d
ExecStop=/usr/bin/docker compose -f /opt/psyched/containers/cerebellum/docker-compose.yml down

[Install]
WantedBy=multi-user.target

