# systemd/layer1-bridge.service
[Unit]
Description=Zenoh <-> ROS2 DDS bridge
# layer1-zenoh is optional; do not require the zenoh fabric at boot
# After and Requires removed to allow the bridge to start independently

[Service]
Type=simple
ExecStart=/usr/bin/zenoh-bridge-ros2dds -c /etc/zenoh/bridge_ros2.json5
Restart=always

[Install]
WantedBy=multi-user.target
