# systemd/psyche-imu.service
[Unit]
Description=Psyche IMU (MPU6050) service
After=network.target

[Service]
Type=simple
User=pete
Group=dialout
WorkingDirectory=/opt/psyched
ExecStart=/usr/local/bin/psyche-ros2-wrapper.sh ros2 launch ros2_mpu6050_driver mpu6050.launch.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
