# systemd/psyched-updater.service
[Unit]
Description=Update psyched repo and re-apply provisioning
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/psyched-update

