# systemd/ros2-camera-repub.service
[Unit]
Description=ROS 2 republisher for Zenoh camera
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Restart=always
Environment=IN_KEY=rt/vision/cam0
Environment=OUT_TOPIC=/camera/cam0/image/compressed
ExecStart=/usr/local/bin/ros_exec /usr/bin/python3 /opt/psyched/layer3/services/ros2_republish_camera.py \
  --in-key ${IN_KEY} --out-topic ${OUT_TOPIC}

[Install]
WantedBy=multi-user.target

