# systemd/ros2-audio-repub.service
[Unit]
Description=ROS 2 republisher for Zenoh audio
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Restart=always
Environment=IN_KEY=rt/audio/ear0
Environment=OUT_TOPIC=/audio/ear0/raw
Environment=AUDIO_RATE=16000
Environment=AUDIO_CHANNELS=1
ExecStart=/usr/local/bin/ros_exec /usr/bin/python3 /opt/psyched/layer3/services/ros2_republish_audio.py \
  --in-key ${IN_KEY} --out-topic ${OUT_TOPIC} --rate ${AUDIO_RATE} --channels ${AUDIO_CHANNELS}

[Install]
WantedBy=multi-user.target

