# systemd/zenoh-gnss-pub.service
[Unit]
Description=Zenoh GNSS NMEA publisher (ear)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pete
Group=dialout
RuntimeDirectory=psyched
Environment=GNSS_DEVICE=/dev/ttyACM0
Environment=GNSS_BAUD=9600
Environment=ZENOH_CONFIG=
ExecStart=/opt/psyched/venv/bin/python -u /opt/psyched/layer3/services/zenoh_gnss_pub.py
Restart=on-failure
TimeoutStartSec=20
# device access managed via dialout group or udev rules ensuring pete ownership

[Install]
WantedBy=multi-user.target
