# systemd/zenoh-audio-pub.service
[Unit]
Description=Zenoh microphone PCM publisher (ear)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Restart=always
Environment=AUDIO_KEY=rt/audio/ear0
Environment=AUDIO_RATE=16000
Environment=AUDIO_CHUNK_MS=20
Environment=AUDIO_CHANNELS=1
Environment=AUDIO_DEVICE=default
ExecStart=/usr/bin/python3 /opt/psyched/layer3/services/zenoh_audio_pub.py

[Install]
WantedBy=multi-user.target

